# docker/s6/services.d/app/run
#!/usr/bin/execlineb -P

nodejs /opt/bulk-sms/bin/www
